# meta-kernel/recipes-kernel/mknod-hello/files/mknod-hello.service
[Unit]
Description=Create a node for the hello char driver

[Service]
ExecStart=/usr/bin/mknod-hello.sh

[Install]
WantedBy=multi-user.target
